# root/etc/crontab
@reboot  /usr/bin/lego-certbot >> /var/log/cron 2>&1
47 11 14 * * /usr/bin/lego-certbot >> /var/log/cron 2>&1